# backend/requests/crontab
* * * * * /usr/local/bin/python /app/autoreject.py >> /var/log/cron.log 2>&1
